# com.onebighing.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onebighing.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/OneBigThingFloat.app/Contents/MacOS/OneBigThingFloat</string>
        <string>--prompt</string>
    </array>

    <!-- Run at 8:00 AM every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Also run when user logs in (catches missed schedules) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep trying if it fails -->
    <key>KeepAlive</key>
    <false/>

    <!-- Standard output/error logging (optional) -->
    <key>StandardOutPath</key>
    <string>/tmp/onebighing.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/onebighing.error.log</string>
</dict>
</plist>
